# msmd/msmd/msmd.service
[Unit]
Description=Monkey Hi Hat TCP Relay

[Service]
Type=notify
ExecStart=/usr/sbin/msmd

[Install]
WantedBy=multi-user.target